# root/etc/s6-overlay/s6-rc.d/svc-webui/run
#!/command/with-contenv bash
exec /vpn-gateway-ui
